# backup-automation/launchd/com.backup.automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.backup.automation</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/guipalm4/Dev/Scripts/backup-automation/src/backup.sh</string>
        <string>incremental</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>16</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/guipalm4/Dev/Scripts/backup-automation/logs/launchd.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/guipalm4/Dev/Scripts/backup-automation/logs/launchd_error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/guipalm4</string>
    </dict>
</dict>
</plist>